# scripts/liotad.service
[Unit]
Description=Liota

[Service]
User=root
ExecStart=/etc/liota/packages/liotad.py

[Install]
WantedBy=multi-user.target

